# deploy/hostinger/vs-dashboard-api.service
[Unit]
Description=VS Dashboard FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/vs-dashboard/backend
EnvironmentFile=/opt/vs-dashboard/backend/.env
ExecStart=/opt/vs-dashboard/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
